# ec2-user@ec2-13-57-223-130.us-west-1.compute.amazonaws.com/ether-fund-react/ec2-user@ec2-13-57-223-130.us-west-1.compute.amazonaws.com/ether-fund-react/ec2-user@ec2-13-57-223-130.us-west-1.compute.amazonaws.com/var/lib/dpkg/info/mongodb-org-server.prerm
#!/bin/sh
set -e
# Automatically added by dh_installinit
if [ -x "/etc/init.d/mongod" ] || [ -e "/etc/init/mongod.conf" ]; then
	invoke-rc.d mongod stop || exit $?
fi
# End automatically added section
