# ec2-user@ec2-13-57-223-130.us-west-1.compute.amazonaws.com/server/ec2-user@ec2-13-57-223-130.us-west-1.compute.amazonaws.com/ether-fund-react/ec2-user@ec2-13-57-223-130.us-west-1.compute.amazonaws.com/var/lib/dpkg/info/slack-desktop.list
/.
/usr
/usr/bin
/usr/share
/usr/share/applications
/usr/share/applications/slack.desktop
/usr/share/pixmaps
/usr/share/pixmaps/slack.png
/usr/share/lintian
/usr/share/lintian/overrides
/usr/share/lintian/overrides/slack
/usr/share/doc
/usr/share/doc/slack-desktop
/usr/share/doc/slack-desktop/copyright
/usr/lib
/usr/lib/slack
/usr/lib/slack/content_resources_200_percent.pak
/usr/lib/slack/LICENSES-linux.json
/usr/lib/slack/resources
/usr/lib/slack/resources/electron.asar
/usr/lib/slack/resources/app.asar.unpacked
/usr/lib/slack/resources/app.asar.unpacked/node_modules
/usr/lib/slack/resources/app.asar.unpacked/node_modules/windows-quiet-hours
/usr/lib/slack/resources/app.asar.unpacked/node_modules/windows-quiet-hours/build
/usr/lib/slack/resources/app.asar.unpacked/node_modules/windows-quiet-hours/build/Release
/usr/lib/slack/resources/app.asar.unpacked/node_modules/windows-quiet-hours/build/Release/quiethours.node
/usr/lib/slack/resources/app.asar.unpacked/node_modules/edge-atom-shell
/usr/lib/slack/resources/app.asar.unpacked/node_modules/edge-atom-shell/build
/usr/lib/slack/resources/app.asar.unpacked/node_modules/edge-atom-shell/build/Release
/usr/lib/slack/resources/app.asar.unpacked/node_modules/edge-atom-shell/build/Release/edge.node
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/edge-cs
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/edge-cs/lib
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/edge-cs/lib/edge-cs.dll
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/cld
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/cld/build
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/cld/build/Release
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/cld/build/Release/cld.node
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/v8-profiler
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/v8-profiler/build
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/v8-profiler/build/Release
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/v8-profiler/build/Release/profiler.node
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/system-idle-time
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/system-idle-time/build
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/system-idle-time/build/Release
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/system-idle-time/build/Release/system_idle_time.node
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/spellchecker
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/spellchecker/build
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/spellchecker/build/Release
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@paulcbetts/spellchecker/build/Release/spellchecker.node
/usr/lib/slack/resources/app.asar.unpacked/node_modules/windows-active-process
/usr/lib/slack/resources/app.asar.unpacked/node_modules/windows-active-process/build
/usr/lib/slack/resources/app.asar.unpacked/node_modules/windows-active-process/build/Release
/usr/lib/slack/resources/app.asar.unpacked/node_modules/windows-active-process/build/Release/processname.node
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@slack
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@slack/macgap-pref-reader
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@slack/macgap-pref-reader/build
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@slack/macgap-pref-reader/build/Release
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@slack/macgap-pref-reader/build/Release/macgap_pref_reader.node
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@slack/slack-calls
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@slack/slack-calls/build
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@slack/slack-calls/build/Release
/usr/lib/slack/resources/app.asar.unpacked/node_modules/@slack/slack-calls/build/Release/slack-calls.node
/usr/lib/slack/resources/app.asar.unpacked/node_modules/runas
/usr/lib/slack/resources/app.asar.unpacked/node_modules/runas/build
/usr/lib/slack/resources/app.asar.unpacked/node_modules/runas/build/Release
/usr/lib/slack/resources/app.asar.unpacked/node_modules/runas/build/Release/runas.node
/usr/lib/slack/resources/app.asar.unpacked/node_modules/nslog
/usr/lib/slack/resources/app.asar.unpacked/node_modules/nslog/build
/usr/lib/slack/resources/app.asar.unpacked/node_modules/nslog/build/Release
/usr/lib/slack/resources/app.asar.unpacked/node_modules/nslog/build/Release/nslog.node
/usr/lib/slack/resources/app.asar.unpacked/node_modules/keyboard-layout
/usr/lib/slack/resources/app.asar.unpacked/node_modules/keyboard-layout/build
/usr/lib/slack/resources/app.asar.unpacked/node_modules/keyboard-layout/build/Release
/usr/lib/slack/resources/app.asar.unpacked/node_modules/keyboard-layout/build/Release/keyboard-layout-manager.node
/usr/lib/slack/resources/app.asar.unpacked/src
/usr/lib/slack/resources/app.asar.unpacked/src/static
/usr/lib/slack/resources/app.asar.unpacked/src/static/lato.less
/usr/lib/slack/resources/app.asar.unpacked/src/static/slack-menubar-rest-Template@2x.png
/usr/lib/slack/resources/app.asar.unpacked/src/static/ts_icon_times_circle.svg
/usr/lib/slack/resources/app.asar.unpacked/src/static/index.js
/usr/lib/slack/resources/app.asar.unpacked/src/static/Hash.png
/usr/lib/slack/resources/app.asar.unpacked/src/static/logo_black@3x.png
/usr/lib/slack/resources/app.asar.unpacked/src/static/ssb-interop.js
/usr/lib/slack/resources/app.asar.unpacked/src/static/slack-menubar-highlight.png
/usr/lib/slack/resources/app.asar.unpacked/src/static/react-virtualized.less
/usr/lib/slack/resources/app.asar.unpacked/src/static/app-win10.ico
/usr/lib/slack/resources/app.asar.unpacked/src/static/slack-menubar-unread-Template@2x.png
/usr/lib/slack/resources/app.asar.unpacked/src/static/slack-menubar-pressed@2x.png
/usr/lib/slack/resources/app.asar.unpacked/src/static/slack-taskbar-unread.png
/usr/lib/slack/resources/app.asar.unpacked/src/static/Hash@2x.png
/usr/lib/slack/resources/app.asar.unpacked/src/static/slack-icons-font.less
/usr/lib/slack/resources/app.asar.unpacked/src/static/bugsnag-3.js
/usr/lib/slack/resources/app.asar.unpacked/src/static/index.jade
/usr/lib/slack/resources/app.asar.unpacked/src/static/slack-menubar-unread-Template.png
/usr/lib/slack/resources/app.asar.unpacked/src/static/slack-menubar-pressed.png
/usr/lib/slack/resources/app.asar.unpacked/src/static/slack-menubar-highlight@2x.png
/usr/lib/slack/resources/app.asar.unpacked/src/static/plastic
/usr/lib/slack/resources/app.asar.unpacked/src/static/plastic/plastic_typography.less
/usr/lib/slack/resources/app.asar.unpacked/src/static/plastic/plastic_colors.less
/usr/lib/slack/resources/app.asar.unpacked/src/static/plastic/plastic_utilities.less
/usr/lib/slack/resources/app.asar.unpacked/src/static/slack-menubar-rest-Template.png
/usr/lib/slack/resources/app.asar.unpacked/src/static/override
/usr/lib/slack/resources/app.asar.unpacked/src/static/override/dropins.js
/usr/lib/slack/resources/app.asar.unpacked/src/static/ssb-interop-lite.js
/usr/lib/slack/resources/app.asar.unpacked/src/static/slack-taskbar-highlight.png
/usr/lib/slack/resources/app.asar.unpacked/src/static/calls-test-sound.pcm
/usr/lib/slack/resources/app.asar.unpacked/src/static/slack.png
/usr/lib/slack/resources/app.asar.unpacked/src/static/Hash@3x.png
/usr/lib/slack/resources/app.asar.unpacked/src/static/logo_black@2x.png
/usr/lib/slack/resources/app.asar.unpacked/src/static/logo_white@2x.png
/usr/lib/slack/resources/app.asar.unpacked/src/static/slack-taskbar-rest.png
/usr/lib/slack/resources/app.asar.unpacked/src/static/ts_icon_arrow_large_left.svg
/usr/lib/slack/resources/app.asar.unpacked/src/static/normalize.css
/usr/lib/slack/resources/app.asar.unpacked/src/static/logo_white@3x.png
/usr/lib/slack/resources/app.asar
/usr/lib/slack/content_shell.pak
/usr/lib/slack/version
/usr/lib/slack/locales
/usr/lib/slack/locales/ca.pak
/usr/lib/slack/locales/nl.pak
/usr/lib/slack/locales/uk.pak
/usr/lib/slack/locales/cs.pak
/usr/lib/slack/locales/ar.pak
/usr/lib/slack/locales/en-US.pak
/usr/lib/slack/locales/fil.pak
/usr/lib/slack/locales/el.pak
/usr/lib/slack/locales/gu.pak
/usr/lib/slack/locales/en-GB.pak
/usr/lib/slack/locales/sk.pak
/usr/lib/slack/locales/bn.pak
/usr/lib/slack/locales/ml.pak
/usr/lib/slack/locales/es.pak
/usr/lib/slack/locales/pt-PT.pak
/usr/lib/slack/locales/lt.pak
/usr/lib/slack/locales/ta.pak
/usr/lib/slack/locales/es-419.pak
/usr/lib/slack/locales/vi.pak
/usr/lib/slack/locales/da.pak
/usr/lib/slack/locales/tr.pak
/usr/lib/slack/locales/zh-CN.pak
/usr/lib/slack/locales/sw.pak
/usr/lib/slack/locales/nb.pak
/usr/lib/slack/locales/th.pak
/usr/lib/slack/locales/ko.pak
/usr/lib/slack/locales/ro.pak
/usr/lib/slack/locales/sr.pak
/usr/lib/slack/locales/mr.pak
/usr/lib/slack/locales/kn.pak
/usr/lib/slack/locales/it.pak
/usr/lib/slack/locales/hi.pak
/usr/lib/slack/locales/fake-bidi.pak
/usr/lib/slack/locales/et.pak
/usr/lib/slack/locales/pt-BR.pak
/usr/lib/slack/locales/hu.pak
/usr/lib/slack/locales/sl.pak
/usr/lib/slack/locales/de.pak
/usr/lib/slack/locales/ja.pak
/usr/lib/slack/locales/id.pak
/usr/lib/slack/locales/hr.pak
/usr/lib/slack/locales/sv.pak
/usr/lib/slack/locales/fr.pak
/usr/lib/slack/locales/pl.pak
/usr/lib/slack/locales/am.pak
/usr/lib/slack/locales/ms.pak
/usr/lib/slack/locales/fa.pak
/usr/lib/slack/locales/te.pak
/usr/lib/slack/locales/fi.pak
/usr/lib/slack/locales/zh-TW.pak
/usr/lib/slack/locales/ru.pak
/usr/lib/slack/locales/lv.pak
/usr/lib/slack/locales/bg.pak
/usr/lib/slack/locales/he.pak
/usr/lib/slack/libffmpeg.so
/usr/lib/slack/blink_image_resources_200_percent.pak
/usr/lib/slack/snapshot_blob.bin
/usr/lib/slack/pdf_viewer_resources.pak
/usr/lib/slack/ui_resources_200_percent.pak
/usr/lib/slack/icudtl.dat
/usr/lib/slack/libnode.so
/usr/lib/slack/slack
/usr/lib/slack/natives_blob.bin
/usr/lib/slack/libCallsCore.so
/usr/lib/slack/views_resources_200_percent.pak
/usr/lib/slack/LICENSE
/etc
/etc/cron.daily
/etc/cron.daily/slack
/usr/bin/slack
