# ec2-user@ec2-13-57-223-130.us-west-1.compute.amazonaws.com/ether-fund-react/ec2-user@ec2-13-57-223-130.us-west-1.compute.amazonaws.com/run/systemd/generator.late/apache-htcacheclean.service
# Automatically generated by systemd-sysv-generator

[Unit]
Documentation=man:systemd-sysv-generator(8)
SourcePath=/etc/init.d/apache-htcacheclean
Description=LSB: Cache cleaner process for Apache2 web server
Before=shutdown.target
After=remote-fs.target
After=systemd-journald-dev-log.socket
Conflicts=shutdown.target

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/etc/init.d/apache-htcacheclean start
ExecStop=/etc/init.d/apache-htcacheclean stop
