# ec2-user@ec2-13-57-223-130.us-west-1.compute.amazonaws.com/run/systemd/generator/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
Before=local-fs.target

[Mount]
What=/dev/disk/by-uuid/3f05c991-6d8c-46f7-a857-4c63a4ff4f07
Where=/
Type=ext4
Options=errors=remount-ro
